# taiga-back
#!/sbin/openrc-run
#------------------------------------------------------------------------------
# Modified Calculate Utilities 3.7.2.17
# Processed template files:
# /var/calculate/templates/taiga/init.taiga
# To modify this file, create a /etc/init.d/taiga.clt template.
#------------------------------------------------------------------------------
# Copyright 2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Taiga daemon"
description=""
pidfile=/run/taiga-back.pid
command_user=taiga
output_log=/var/log/taiga-back.log
error_log=/var/log/taiga-back.log
directory=/var/calculate/www/taiga/taiga-back
start_stop_daemon_args=""
command="/var/calculate/www/taiga/taiga-back/.venv/bin/gunicorn"
command_args="--workers 4 --timeout 60 --log-level=info --access-logfile - --bind 0.0.0.0:8001 taiga.wsgi"
command_background=true
source /var/calculate/www/taiga/taiga-back/.venv/bin/activate

depend() {
    need nginx rabbitmq postgresql
    need taiga-async taiga-events taiga-protected
}

start_pre() {
    checkpath -f -o taiga -m 0600 /var/log/taiga-back.log
}
